# rc.local
#!/bin/sh
#exemplo de inicializacao para inserir sensores i2c do barramento principal do ev3

beep -f 200
sh -c 'echo 0 > /sys/class/leds/led0\:green\:brick-status/brightness'
sh -c 'echo 0 > /sys/class/leds/led1\:green\:brick-status/brightness'

/home/robot/servos.sh

beep -f 600
sh -c 'echo 255 > /sys/class/leds/led0\:green\:brick-status/brightness'
sh -c 'echo 255 > /sys/class/leds/led1\:green\:brick-status/brightness'

exit 0
